# deploy/webhook-server.service
[Unit]
Description=Marxnager Bot Webhook Auto-Deploy Server
Documentation=https://github.com/your-username/Sindicato-tg-bot-1
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt
Environment="PORT=5000"
Environment="WEBHOOK_SECRET=your-github-webhook-secret"
Environment="UPDATE_SCRIPT=/usr/local/bin/update-marxnager.sh"
Environment="LOG_FILE=/var/log/webhook-server.log"
ExecStart=/usr/bin/python3 /opt/webhook-server.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
